# vendor/templates/puma.daemon.erb
#!/bin/bash

# <%= "Start or stop #{fetch(:application)}-puma" %>
#
# <%= fetch(:daemon_author_name) %>

### BEGIN INIT INFO
# Provides:          <%= "#{fetch(:application)}-puma" %>
# Required-Start:    $local_fs $remote_fs $syslog $network
# Required-Stop:     $local_fs $remote_fs $syslog $network
# Should-Start:
# Should-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: <%= "Starts #{fetch(:application)}-puma" %>
# Description:       <%= "Launches #{fetch(:application)}-puma" %>
### END INIT INFO

PATH=/home/<%=fetch(:monit_user)%>/.rbenv/shims:/home/<%=fetch(:monit_user)%>/.rbenv/bin:$PATH
PIDFILE=<%= "#{shared_path}/pids/puma.pid" %>
export RAILS_ENV=production

getPID() {
  if [ ! -z `cat $PIDFILE 2> /dev/null` ]; then
    if [[ ! -z $(ps -ef | grep -w `cat $PIDFILE` | grep puma | grep -v grep) ]]; then
      echo `cat $PIDFILE 2> /dev/null`
    fi
  fi
}

case "$1" in
  start)
    echo -ne "<%= "Starting #{fetch(:application)}-puma..." %>"

    if [ -z "$(getPID)" ]; then
      cd <%=current_path%>
      bundle exec puma -q -d -C <%= "#{shared_path}/config/puma.rb" %>

      echo "started!"
    else
      echo "is already running!"
    fi
    ;;

  stop)
    echo -ne "<%= "Stopping #{fetch(:application)}-puma..." %>"

    if [ -z "$(getPID)" ]; then
      echo "is not running!"
    else
      cd <%=current_path%>
      bundle exec pumactl -S <%= "#{shared_path}/pids/puma.state" %> stop

      echo "stopped!"
    fi
    ;;

  restart)
    $0 stop
    sleep 3
    $0 start
    ;;

  check_pid)
    if [ -z "$(getPID)" ]; then
      > $PIDFILE
    fi
    ;;

  reload)
    echo -ne "<%= "Reloading #{fetch(:application)}puma..." %>"

    if [ -z "$(getPID)" ]; then
      echo "is not running!"
    else
      kill -USR2 $(getPID) &> /dev/null

      echo "reloaded!"
    fi
    ;;

  *)
    echo "usage: $0 {start|stop|restart}"
esac
exit 0
